# deploy/systemd/linofinance-api.service
[Unit]
Description=LinoFinance API
After=network.target

[Service]
Type=simple
User=linofinance
Group=linofinance
WorkingDirectory=/srv/linofinance/backend
EnvironmentFile=/etc/linofinance/api.env
ExecStart=/srv/linofinance/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
