# deploy/systemd/linofinance-jobs.service
[Unit]
Description=LinoFinance scheduled jobs (credit-due reminders + attachment cleanup)
After=network.target linofinance-api.service

[Service]
Type=oneshot
User=linofinance
Group=linofinance
WorkingDirectory=/opt/linofinance/app/current/backend
EnvironmentFile=/etc/linofinance/api.env
ExecStart=/opt/linofinance/app/current/backend/.venv/bin/python scripts/run_scheduled_jobs.py
